# docker/openresty/root/var/spool/cron/crontabs/root

*/6 * * * * /usr/sbin/logrotate /etc/logrotate.d/interaction-logger.conf
# leave an empty line

