# docker/openresty/root/etc/services.d/openresty/run
#!/bin/bash
exec openresty -g "daemon off;"